# codedeploy-scripts/upstart/go-rest-api-template.conf
# upstart configuration

description "golang server"
author "lee provoost"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
exec /home/ec2-user/go/bin/go-rest-api-template
